# DEBIAN/postinst
#!/usr/bin/env sh
cd /etc/dischord
dischord -wizard -quit_after_wizard
cd -
systemctl enable --now dischord.service
