# DEBIAN/prerm
#!/usr/bin/env sh
systemctl disable --now dischord.service || true
